# controlRoom_server/server/crontab/crontab_10302019

# CAO
# 0 11 * * * /home/hnpcen/bin/shell/gold_cao_11.sh 2 > /home/hnpcen/bin/shell/gold_cao.sh_ma_11.log
# 0 14 * * * /home/hnpcen/bin/shell/gold_cao_14.sh 2 >  /home/hnpcen/bin/shell/gold_cao.sh_ma_14.log
# 0 7 * * * . /home/hnpcen/env/envCEN; /home/hnpcen/bin/shell/loadPromoImpactCalculationForCGO/loadPromoImpactCalc.sh > /home/hnpcen/bin/shell/loadPromoImpactCalculationForCGO/loadPromoImpactCalc.log

#Nightly Restart GMS web service hnpcen
#15 00 * * * /home/hnpcen/bin/shell/restart_processes.sh > /home/hnpcen/bin/shell/restart_processes.log 2>&1

# ICR
*/15 3-15 * * * /opt/apps/controlRoom/controlRoom_server/server/scripts/EDI_invoice_alert.sh
0 19 * * 7 /opt/apps/controlRoom/controlRoom_server/server/scripts/EDI_invoice_WEEKLY_alert.sh
30 11 * * * /opt/apps/controlRoom/controlRoom_server/server/scripts/CGO_negative_inventory.sh
30 11 * * * /opt/apps/controlRoom/controlRoom_server/server/scripts/CGO_Orders_Recap.sh

0 * * * * /opt/apps/controlRoom/controlRoom_server/server/start_backend.sh
0 * * * * /opt/apps/controlRoom/controlRoom_batch_server/server/start_batch_backend.sh
