# controlRoom_server/server/crontab/crontab_prod_10302019
#00 2 * * * $HOME/bin/shell/restart_TomEE_hnpcen.sh
01 0 * * * $HOME/seyc/nmon/start_nmon.sh
0 1 * * * /home/hnpcen/bin/deleteFiles /data/hnpcen/gaia_out txt.2* 1
1 1 * * * /home/hnpcen/bin/deleteFiles /data/hnpcen/gaia_out log 1
2 1 * * * /home/hnpcen/bin/deleteFiles /data/hnpcen/tomee_log log 14
3 1 * * * /home/hnpcen/bin/deleteFiles /data/hnpcen/tomee_log txt 14
4 1 * * * /home/hnpcen/bin/deleteFiles /data/hnpcen/tomee_log gz 14
5 1 * * * /home/hnpcen/bin/deleteFiles /data/hnpcen/bin_log log 14
6 1 * * * /home/hnpcen/bin/deleteFiles /data/hnpcen/nmon_files nmon 365
7 1 * * * /home/hnpcen/bin/deleteFiles /data/hnpcen/nmon_files gz 365
8 1 * * * /home/hnpcen/bin/deleteFiles /data/hnpcen/ws_log gz 14
9 1 * * * /home/hnpcen/bin/deleteFiles /data/hnpcen/ws_log txt 14
#0 5 * * * /home/hnpcen/bin/shell/gold_cao.sh 2 > /home/hnpcen/bin/shell/gold_cao.sh_ma.log
0 11 * * * /home/hnpcen/bin/shell/gold_cao_11.sh 2 > /home/hnpcen/bin/shell/gold_cao.sh_ma_11.log
0 14 * * * /home/hnpcen/bin/shell/gold_cao_14.sh 2 >  /home/hnpcen/bin/shell/gold_cao.sh_ma_14.log
0 7 * * * . /home/hnpcen/env/envCEN; /home/hnpcen/bin/shell/loadPromoImpactCalculationForCGO/loadPromoImpactCalc.sh > /home/hnpcen/bin/shell/loadPromoImpactCalculationForCGO/loadPromoImpactCalc.log
0 6 * * * /home/hnpcen/bin/shell/abe/pi_neg_inventory.sh

#Nightly Restart GMS web service hnpcen
#15 00 * * * /home/hnpcen/bin/shell/restart_processes.sh > /home/hnpcen/bin/shell/restart_processes.log 2>&1

# PICS/RGIS automatic integration counting
#*/15 1-8 * * * /home/hnpcen/bin/shell/thirdpartyinventory/0.GetFTPfiles.sh
